# .config/autostart/mate-polkit.desktop
[Desktop Entry]
Type=Application
Name=MATE Polkit Agent
Exec=/usr/lib/mate-polkit/polkit-mate-authentication-agent-1


