# .config/autostart/xsettings.desktop
[Desktop Entry]
Type=Application
Name=XSettings-Daemon
Exec=systemctl --user start xsettingsd

