# services/nm-applet.service
[Unit]
Description=NetworkManager Applet

[Service]
Type=simple
ExecStart=/usr/bin/nm-applet
Restart=on-failure

