# services/waybar.service
[Unit]
Description=Waybar (Wayland Status Bar)

[Service]
Type=simple
ExecStart=/usr/bin/waybar
Restart=on-failure


